# sys_files/etc/profile.d/00-aliases.sh
function atomic() {
    just -f /etc/Justfile $@
}
