# sys_files/etc/systemd/user/bootstrap.service
[Unit]
Description=Bootstrap installation

[Service]
ExecStart=/usr/local/bin/__post-install

[Install]
RequiredBy=default.target
