# sys_files/etc/systemd/user/rotate-wallpaper.timer
[Unit]
Description=Rotate the current `swww` wallpaper with spotligth wallpapers
Wants=graphical.target

[Timer]
OnBootSec=20min
OnUnitActiveSec=120min
Unit=rotate-wallpaper.service

[Install]
WantedBy=timers.target
